# config/captive-portal.service
# config/captive-portal.service
# Copier vers /etc/systemd/system/captive-portal.service
# sudo systemctl enable captive-portal
# sudo systemctl start captive-portal

[Unit]
Description=Captive Portal SSO
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/captive-portal

# Environnement
EnvironmentFile=/opt/captive-portal/.env

# Démarrage
ExecStart=/opt/captive-portal/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 4

# Redémarrage automatique
Restart=always
RestartSec=5

# Limites
LimitNOFILE=65535

# Capabilities pour nftables
AmbientCapabilities=CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_ADMIN

# Sécurité
NoNewPrivileges=false
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/captive-portal

[Install]
WantedBy=multi-user.target
